# rootfs/etc/services.d/logmonitor/run
#!/bin/execlineb

# Exit and disable service if no file to be monitored.
ifelse -n
{
  pipeline { cat /etc/logmonitor/logmonitor.conf }
  s6-grep -q "^[ \t]*LOG_FILES=/.\\+"
}
{
  foreground { s6-echo "[logmonitor] no file to monitor: disabling service..." }
  foreground { s6-touch /var/run/s6/services/logmonitor/down }
  foreground { s6-svc -xd /var/run/s6/services/logmonitor }
  exit 0
}

# Get log files to be monitored from configuration.
backtick -I -n LOG_FILES
{
  pipeline { cat /etc/logmonitor/logmonitor.conf }
  pipeline { s6-grep "^[ \t]*LOG_FILES=" }
  s6-cut -d "=" -f2-
}
importas -u -sd "," LOG_FILES LOG_FILES

# Start the log monitor.
foreground { s6-echo "[logmonitor] starting..." }
exec /usr/bin/logmonitor ${LOG_FILES}

# vim: set ft=sh :
